# Red Hat System Administration II/Lab-09-Systemd-Timers-For-Task-Scheduling/maintenance.timer
[Unit]
Description=Run maintenance daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=maintenance.service

[Install]
WantedBy=timers.target
